# systemd/mock-fw.socket
[Unit]
Description=Mock FBNIC Firmware Socket

[Socket]
ListenStream=/tmp/fbnic-ctrl-skt
Accept=yes

[Install]
WantedBy=sockets.target
