# systemd/mock-fw@.service
[Unit]
Description=Mock FBNIC Firmware (connection %i)

[Service]
ExecStart=/path/to/mock_fw.py
